# rtcm-stream-gateway/deploy.service.example
[Unit]
Description=RTCM stream gateway
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=rtcmgw
Group=rtcmgw
WorkingDirectory=/opt/rtcm-stream-gateway
EnvironmentFile=/opt/rtcm-stream-gateway/gateway.env
ExecStart=/opt/rtcm-stream-gateway/bin/gateway
Restart=always
RestartSec=2
AmbientCapabilities=CAP_NET_RAW CAP_NET_ADMIN
CapabilityBoundingSet=CAP_NET_RAW CAP_NET_ADMIN
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=full
ProtectHome=true
ReadWritePaths=/opt/rtcm-stream-gateway
CPUQuota=50%
MemoryMax=900M

[Install]
WantedBy=multi-user.target
